# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true

[program:php-fpm5.6]
command=/usr/sbin/php-fpm5.6 -F
autorestart=true